# conf/value.conf
[program:value]
command = /webapps/value_tool/bin/gunicorn_start                      ; Command to start app
user = value                                                          ; User to run as
stdout_logfile = /webapps/value_tool/logs/gunicorn_supervisor.log     ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8